# files/mhnd_bootstrap.conf
description "Development Machine Bootstrap"

start on (starting ssh)

console output

env __HOST

exec /bin/bash /mhnd/bootstrap.sh
